# services/comp8.service
[Unit]
Description=Competition 8 — All 5 Agents Max Profit 20 pairs no limit (port 8129)
After=network.target candle.service
Wants=candle.service

[Service]
Type=simple
User=mrpolo
WorkingDirectory=/home/mrpolo/binanceBot
ExecStart=/usr/bin/python3 -m paper8.comp8_server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
